# deploy/systemd/sage-wiki-bridge.service
[Unit]
Description=sage-wiki WeChat Official Account bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=sagewiki
Group=sagewiki
WorkingDirectory=/data/workspace/sage-wiki-bridge-wxo
EnvironmentFile=/data/workspace/sage-wiki-bridge-wxo/.env
Environment=ENV_FILE=/data/workspace/sage-wiki-bridge-wxo/.env
ExecStart=/data/workspace/sage-wiki-bridge-wxo/scripts/bridgectl.sh run
Restart=on-failure
RestartSec=5s
TimeoutStopSec=20s

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/data/workspace/sage-wiki-bridge-wxo/data /data/workspace/sage-wiki/source

MemoryAccounting=true
MemoryMax=256M

[Install]
WantedBy=multi-user.target
